# install/cloud/gce/systemd/cockroach.service.d/30-envvars.conf
[Service]
# Needed to workaround issue - cockroachdb/cockroach#12675
Environment="COCKROACH_METRICS_SAMPLE_INTERVAL=1000h"
